# cmd/steamos/conf/zaparoo.service
[Unit]
Description=Zaparoo Core service

[Service]
Type=exec
Restart=on-failure
RestartSec=5
StandardError=syslog
User=deck
Group=deck
WorkingDirectory=%%WORKING%%
ExecStart=%%EXEC%%

[Install]
WantedBy=multi-user.target